# cron.hourly/motd
#!/bin/sh

pacman -Sy
echo -e "$(run-parts /etc/motd.d)\n" > /etc/motd
