# cron.hourly/motd
#!/bin/sh

echo -e "$(run-parts /etc/motd.d)\n" > /etc/motd
